# DEBIAN/postinst
if [ -e /bin/systemctl ]; then
	systemctl enable rpi_temp_monitor.service
fi
